# automatico.cron
#!/bin/bash
tmux new-session -d -s biomercados
tmux send-keys -t biomercados "cd /var/www/html/e-commerce-bio" C-m
tmux send-keys -t biomercados "php iniciar_servidor.php" C-m
